# deploy/node-agent.service
# systemd unit for llm-router node agent (sandboxed)
# Install: sudo cp node-agent.service /etc/systemd/system/llm-router-agent.service
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now llm-router-agent
#
# Requires: useradd --system --shell /usr/sbin/nologin llm-router

[Unit]
Description=LLM Router Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llm-router
Group=llm-router

ExecStart=/home/erewhon/Projects/erewhon/llm-router/.venv/bin/python -m llm_router.node_agent.main --port 8100

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PATH=/home/erewhon/.local/bin:/usr/local/bin:/usr/bin:/bin

# --- State directories (systemd creates & owns these) ---
# /var/lib/llm-router — persistent state (logs, env files for vLLM units)
StateDirectory=llm-router
# /run/llm-router — ephemeral runtime data (PID files)
RuntimeDirectory=llm-router

# --- Filesystem sandboxing ---

# Read-only access to /usr, /boot, /etc; writable only to StateDirectory/RuntimeDirectory
ProtectSystem=strict
# Node agent writes vLLM env files here before starting vllm@ units
ReadWritePaths=/etc/llm-router/vllm-env

# Empty /home, then bind-mount only what the agent needs
ProtectHome=tmpfs
BindReadOnlyPaths=/home/erewhon/Projects/erewhon/llm-router/.venv
BindReadOnlyPaths=/home/erewhon/Projects/erewhon/llm-router/src
BindReadOnlyPaths=/home/erewhon/Projects/erewhon/llm-router/models.yaml
# uv-managed Python (all AI nodes have uv installed)
BindReadOnlyPaths=/home/erewhon/.local/share/uv/python

# Private /tmp and /var/tmp
PrivateTmp=yes

# --- Kernel protections ---
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# --- Capability restrictions ---
# NoNewPrivileges must be off: agent uses sudo to manage vllm@ systemd units
NoNewPrivileges=no
# sudo needs CAP_SETUID/CAP_SETGID to escalate for systemctl
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_READ_SEARCH
AmbientCapabilities=

# --- System call restrictions ---
# No syscall filter: agent needs sudo to manage vllm@ systemd units,
# which requires a broad set of syscalls. Other sandboxing layers
# (ProtectSystem, network isolation, filesystem) provide defense.
SystemCallArchitectures=native

# --- Namespace and address restrictions ---
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
# SUID needed for sudo to manage vllm@ units
RestrictSUIDSGID=no

# --- Device access ---
# GPU sysfs monitoring needs device access (reads /sys/class/drm, runs nvidia-smi)
PrivateDevices=no
# Bind GPU sysfs paths for monitoring
BindReadOnlyPaths=/sys/class/drm

# --- Memory ---
MemoryDenyWriteExecute=no
LockPersonality=yes
UMask=0077

# --- Network isolation ---
# Only allow localhost, LAN (node-to-node communication), deny everything else
IPAddressAllow=127.0.0.0/8 ::1/128 192.168.42.0/24 172.16.0.0/12 100.64.0.0/10
IPAddressDeny=any

[Install]
WantedBy=multi-user.target
