# deploy/qwen3-vl-32b.service
# systemd unit for Qwen3-VL-32B-Thinking (llama.cpp CUDA)
# Deploy to: archimedes
# Install: sudo cp qwen3-vl-32b.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now qwen3-vl-32b

[Unit]
Description=Qwen3-VL-32B-Thinking Vision Model (llama.cpp CUDA)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=erewhon
Group=erewhon
Environment=LD_LIBRARY_PATH=/opt/llama-cpp-cuda:/usr/local/cuda/lib64

ExecStart=/opt/llama-cpp-cuda/llama-server \
    -m /home/erewhon/models/qwen3-vl-32b/Qwen3VL-32B-Thinking-Q4_K_M.gguf \
    --mmproj /home/erewhon/models/qwen3-vl-32b/mmproj-Qwen3VL-32B-Thinking-Q8_0.gguf \
    --host 0.0.0.0 \
    --port 5404 \
    -ngl 999 \
    --ctx-size 32768 \
    --no-webui

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
