# deploy/vllm@.service
# systemd template unit for vLLM inference backends (Docker)
# Install: sudo cp vllm@.service /etc/systemd/system/vllm@.service
#          sudo systemctl daemon-reload
#
# Usage: systemctl start vllm@model-id  (model-id has / replaced with -)
# Config: /etc/llm-router/vllm-env/<instance>.args  (one-line vLLM CLI args)
#
# Requires: Docker with NVIDIA container toolkit
#           mkdir -p /etc/llm-router/vllm-env && chown llm-router /etc/llm-router/vllm-env
#           mkdir -p /opt/llm-vllm/hf-cache  (pre-downloaded models)

[Unit]
Description=vLLM Inference Server (%i)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple

# Default Docker image — override via /etc/llm-router/vllm-env/<instance>.image
Environment=VLLM_IMAGE=vllm/vllm-openai:latest

# Node agent writes vLLM CLI args to a plain text file.
# Optional per-instance image override file (e.g. for nightly).
ExecStartPre=/bin/bash -c 'docker rm -f vllm-%i 2>/dev/null || true'
ExecStart=/bin/bash -c '\
    IMAGE=${VLLM_IMAGE}; \
    [ -f /etc/llm-router/vllm-env/%i.image ] && IMAGE=$(cat /etc/llm-router/vllm-env/%i.image); \
    exec docker run --rm \
        --name vllm-%i \
        --gpus all \
        --network host \
        --ipc host \
        -v /opt/llm-vllm/hf-cache:/root/.cache/huggingface:ro \
        -e HF_HUB_OFFLINE=1 \
        "$IMAGE" \
        $(cat /etc/llm-router/vllm-env/%i.args)'

ExecStop=/usr/bin/docker stop -t 30 vllm-%i

Restart=no
TimeoutStopSec=45
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm@%i

# --- Network isolation ---
# Docker handles network; systemd IPAddressDeny doesn't work with --network=host containers
# The container inherits host networking but vLLM only binds to specified port

[Install]
WantedBy=multi-user.target
